# config/launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.flowclaw.workflow-executor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>workflow-executor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/flowclaw/src</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/flowclaw/logs/workflow-executor.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/flowclaw/logs/workflow-executor-error.log</string>
</dict>
</plist>
