# config/workflow-executor.service
# Platform: Linux ONLY
# For macOS, use config/launchagent.plist instead

[Unit]
Description=FlowClaw Workflow Executor
After=network.target

[Service]
Type=simple
User=your-user
WorkingDirectory=/path/to/flowclaw/src
# Production startup via Gunicorn multi-worker (see src/scripts/start-workflow-executor.sh)
ExecStart=/bin/bash /path/to/flowclaw/src/scripts/start-workflow-executor.sh
Restart=always
RestartSec=5
EnvironmentFile=/path/to/flowclaw/.env

[Install]
WantedBy=multi-user.target
